# litellm.service
[Unit]
Description=LiteLLM Proxy for AI Hustle Machine
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/aimoneymachine
ExecStart=/opt/aimoneymachine/venv/bin/litellm --port 4000 --num_callbacks 0 --drop_params
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
